# docker/daphne.conf
[program:daphne]
command=daphne %(ENV_ASGI_APP)s -u /tmp/nginx.sock.%(process_num)d -v 2
process_name=%(program_name)s-%(process_num)d
numprocs=%(ENV_WEB_PROCESSES)s
killasgroup=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
