# scripts/celerybeat.conf
[program:celerybeat]
command=celery beat -A abas --loglevel=INFO
numprocs=1
killasgroup=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
